# nvr/pi/cron.d/seccam
# NVR napi videók mentése

0 0 * * *     root   /usr/local/bin/seccam
*/5 * * * *   root   /usr/local/bin/seccam_serv