# sbackup/etc/cron.weekly/sbackup
#!/bin/sh
#
#

/usr/local/bin/sbackup-log

#
